# open-falcon/supervisord.conf

[supervisord]
childlogdir = /home/open-falcon/logs/
logfile = /home/open-falcon/logs/supervisord.log
pidfile = /var/run/supervisord.pid
nodaemon = true
directory = /home/open-falcon

[unix_http_server]
file = /var/run/supervisor.sock
chmod = 0777
chown= nobody:nogroup

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///var/run/supervisor.sock

[program:startup]
command = bash /home/open-falcon/scripts/startup.sh
autorestart = false
priority = 1


[program:transfer]
user = work
directory = /home/open-falcon/%(program_name)s
command = /home/open-falcon/%(program_name)s/falcon-%(program_name)s -c /home/open-falcon/cfg/%(program_name)s.cfg.json
stdout_logfile = /home/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/open-falcon/logs/%(program_name)s.err
autostart=false

[program:agent]
user = work
directory = /home/open-falcon/%(program_name)s
command = /home/open-falcon/%(program_name)s/falcon-%(program_name)s -c /home/open-falcon/cfg/%(program_name)s.cfg.json
stdout_logfile = /home/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/open-falcon/logs/%(program_name)s.err
autostart=false

[program:graph]
user = work
directory = /home/open-falcon/%(program_name)s
command = /home/open-falcon/%(program_name)s/falcon-%(program_name)s -c /home/open-falcon/cfg/%(program_name)s.cfg.json
stdout_logfile = /home/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/open-falcon/logs/%(program_name)s.err
autostart=false

[program:query]
user = work
directory = /home/open-falcon/%(program_name)s
command = /home/open-falcon/%(program_name)s/falcon-%(program_name)s -c /home/open-falcon/cfg/%(program_name)s.cfg.json
stdout_logfile = /home/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/open-falcon/logs/%(program_name)s.err
autostart=false

[program:dashboard]
user = work
directory = /home/open-falcon/%(program_name)s
command = /home/open-falcon/%(program_name)s/env/bin/gunicorn -c gunicorn.conf wsgi:app -t 6000 
stdout_logfile = /home/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/open-falcon/logs/%(program_name)s.err
autostart=false

[program:sender]
user = work
directory = /home/open-falcon/%(program_name)s
command = /home/open-falcon/%(program_name)s/falcon-%(program_name)s -c /home/open-falcon/cfg/%(program_name)s.cfg.json
stdout_logfile = /home/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/open-falcon/logs/%(program_name)s.err
autostart=false

[program:fe]
user = work
directory = /home/open-falcon/%(program_name)s
command = /home/open-falcon/%(program_name)s/falcon-%(program_name)s -c /home/open-falcon/cfg/%(program_name)s.cfg.json
stdout_logfile = /home/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/open-falcon/logs/%(program_name)s.err
autostart=false

[program:portal]
user = work
directory = /home/open-falcon/%(program_name)s
command = /home/open-falcon/%(program_name)s/env/bin/gunicorn -c gunicorn.conf wsgi:app -t 6000 
stdout_logfile = /home/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/open-falcon/logs/%(program_name)s.err
autostart=false

[program:hbs]
user = work
directory = /home/open-falcon/%(program_name)s
command = /home/open-falcon/%(program_name)s/falcon-%(program_name)s -c /home/open-falcon/cfg/%(program_name)s.cfg.json
stdout_logfile = /home/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/open-falcon/logs/%(program_name)s.err
autostart=false

[program:judge]
user = work
directory = /home/open-falcon/%(program_name)s
command = /home/open-falcon/%(program_name)s/falcon-%(program_name)s -c /home/open-falcon/cfg/%(program_name)s.cfg.json
stdout_logfile = /home/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/open-falcon/logs/%(program_name)s.err
autostart=false

[program:alarm]
user = work
directory = /home/open-falcon/%(program_name)s
command = /home/open-falcon/%(program_name)s/falcon-%(program_name)s -c /home/open-falcon/cfg/%(program_name)s.cfg.json
stdout_logfile = /home/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/open-falcon/logs/%(program_name)s.err
autostart=false

[program:links]
user = work
directory = /home/open-falcon/%(program_name)s
command = /home/open-falcon/%(program_name)s/env/bin/gunicorn -c gunicorn.conf wsgi:app -t 6000 
stdout_logfile = /home/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/open-falcon/logs/%(program_name)s.err
autostart=false